# roles/unbound/files/unbound-filters.timer
[Unit]
Description=Update Unbound filters

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
